# systemd/home-sync-watcher.service
[Unit]
Description=Home Sync Connectivity Watcher
After=network.target

[Service]
EnvironmentFile=-%h/.config/home-sync/watcher.env
ExecStart=%h/scripts/home-sync-watcher.sh
Restart=always
RestartSec=10

[Install]
WantedBy=default.target
