# systemd/system-services/home-sync-import-accounts@.path
[Unit]
Description=Monitor accounts.json for user %i
After=network.target

[Path]
PathChanged=/home/%i/.config/home-sync/accounts.json

[Install]
WantedBy=multi-user.target
