# systemd/system-services/home-sync-export-accounts.timer
[Unit]
Description=Timer to Export Home-Sync User Accounts

[Timer]
OnCalendar=hourly
Persistent=true

[Install]
WantedBy=timers.target
